# scripts/av/obs/create-obs-monitor-sink.service
[Unit]
Description=Create OBS Monitor Sink (OBS_Monitor)
# PipeWire 세트가 먼저 올라오도록 의존성/순서 강화
After=pipewire.service pipewire-pulse.service wireplumber.service default.target
Wants=pipewire.service pipewire-pulse.service wireplumber.service

[Service]
Type=oneshot
# user 세션/버스/데몬 준비성 확인
ExecStartPre=/usr/bin/env bash -lc 'pactl info >/dev/null'
ExecStart=/usr/local/bin/create-obs-monitor-sink.sh
RemainAfterExit=yes

[Install]
WantedBy=default.target
