# scripts/av/obs/create-obs-virtual-mic.service
[Unit]
Description=Create OBS Virtual Mic (OBS_VirtualMic from OBS_Monitor.monitor)
After=pipewire.service wireplumber.service create-obs-monitor-sink.service default.target
Wants=pipewire.service wireplumber.service create-obs-monitor-sink.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/create-obs-virtual-mic.sh
RemainAfterExit=yes

[Install]
WantedBy=default.target
